# overlay-root/etc/systemd/system/age-tpm-identity.service
[Unit]
Description=Ensure age-plugin-tpm identity exists

[Service]
Type=oneshot
ExecStart=/usr/local/bin/age-tpm-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
