# overlay-root/etc/systemd/system/alertmanager-generate-config.service
[Unit]
Description=Generate Alertmanager config from template with Pushover secrets
Before=alertmanager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/alertmanager-generate-config.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
