# overlay-root/etc/systemd/system/caddy-rootless-preflight.service
[Unit]
Description=Record Caddy rootless migration preflight
After=caddy.service ensure-nas-caddy-account.service sops-distribute-secrets.service systemd-sysctl.service
Requires=caddy.service ensure-nas-caddy-account.service sops-distribute-secrets.service
ConditionPathExists=!/var/lib/nas-migrations/caddy-rootless-preflight-v1/complete

[Service]
Type=oneshot
StateDirectory=nas-migrations
StateDirectoryMode=0700
ExecStart=/usr/local/bin/caddy-rootless-preflight.sh
TimeoutStartSec=infinity
