# overlay-root/etc/systemd/system/caddy-rootless-preflight.timer
[Unit]
Description=Run Caddy rootless migration preflight after boot

[Timer]
OnBootSec=5min
Unit=caddy-rootless-preflight.service
AccuracySec=1min

[Install]
WantedBy=timers.target
