# overlay-root/etc/systemd/system/disk-health-metrics.timer
[Unit]
Description=Collect SMART and ZFS health metrics every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
RandomizedDelaySec=30s

[Install]
WantedBy=timers.target
