# overlay-root/etc/systemd/system/ensure-nas-caddy-account.service
# GENERATED by generate-quadlets.py from quadlets/caddy.toml — DO NOT EDIT
[Unit]
Description=Prepare _nas_caddy for a lingering user manager
DefaultDependencies=no
After=local-fs.target systemd-sysusers.service systemd-tmpfiles-setup.service
Before=sysinit.target systemd-logind.service user@51310.service
ConditionPathExists=/etc/passwd

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ensure-nas-caddy-account.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
