# overlay-root/etc/systemd/system/ensure-nas-grafana-account.service
[Unit]
Description=Ensure _nas_grafana can run a lingering user manager
DefaultDependencies=no
After=local-fs.target systemd-sysusers.service systemd-tmpfiles-setup.service
Before=sysinit.target systemd-logind.service user@51210.service
ConditionPathExists=/etc/passwd

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ensure-nas-grafana-account.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
