# overlay-root/etc/systemd/system/garage-generate-secrets.service
[Unit]
Description=Generate Garage API tokens
After=age-tpm-identity.service
Requires=age-tpm-identity.service
Before=garage.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/garage-generate-secrets.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
