# overlay-root/etc/systemd/system/zfs-health-check.service
[Unit]
Description=ZFS Pool Health Check
Documentation=man:zpool(8)
Requires=zfs.target
After=zfs.target
ConditionACPower=true
ConditionPathIsDirectory=/sys/module/zfs

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zfs-health-check.sh
User=root
StandardOutput=journal
StandardError=journal
PrivateNetwork=yes
ProtectSystem=strict
ReadWritePaths=/dev /sys /proc
RemainAfterExit=no