# docs/templates/rootless-quadlet/ensure-service-account.service.template
[Unit]
Description=Ensure HOST_USER can run a lingering user manager
DefaultDependencies=no
After=local-fs.target systemd-sysusers.service systemd-tmpfiles-setup.service
Before=sysinit.target systemd-logind.service user@HOST_UID.service
ConditionPathExists=/etc/passwd

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ensure-SERVICE_NAME-account.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
