# overlay-root/etc/containers/systemd/garage.container
[Unit]
Description=Garage S3-compatible Object Storage
Documentation=https://garagehq.deuxfleurs.fr/
Wants=network-online.target nss-lookup.target
After=network-online.target nss-lookup.target sops-distribute-secrets.service zfs-create-garage-datasets.service
Requires=sops-distribute-secrets.service zfs-create-garage-datasets.service

[Container]
ContainerName=garage
Image=docker.io/dxflrs/garage:v2.2.0
AutoUpdate=registry

Environment=GARAGE_LOG_TO_JOURNALD=1

# Mount journald socket for logs
Volume=/run/systemd/journal/socket:/run/systemd/journal/socket

# Mount configuration (read-only)
Volume=/etc/garage/garage.toml:/etc/garage.toml:ro,Z

# Secrets injected via podman secrets. The shell driver stores them encrypted
# at rest with systemd-creds in /var/lib/podman-secrets.
Secret=garage-rpc-secret,mode=0400
Secret=garage-admin-token,mode=0400
Secret=garage-metrics-token,mode=0400

# Mount ZFS datasets for data storage (no :Z — datasets are too large for
# recursive SELinux relabeling on every start; labeled at creation time instead)
Volume=/var/lib/garage/meta:/var/lib/garage/meta
Volume=/var/lib/garage/data:/var/lib/garage/data

# Expose ports
# 3900: S3 API
# 3901: RPC (not needed externally for single-node)
# 3902: Web (static site hosting)
# 3903: Admin API
PublishPort=127.0.0.1:3900:3900
PublishPort=127.0.0.1:3902:3902
PublishPort=127.0.0.1:3903:3903

Pull=newer

[Service]
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target default.target
