# overlay-root/etc/profile.d/50_garage.sh
# Set an alias for `garage` which execcutes inside the container
alias garage="sudo podman exec -it garage /garage"